# net_services/plugins/ssh-server/files/etc/uci-defaults/99-prevent-extra-command-error
#!/bin/sh

sed -i '/extra_command "killclients"/d' /etc/init.d/dropbear

exit 0
