# net_core/plugins/ethernet-manager/files/etc/uci-defaults/99-remove-Mxl-scripts
#!/bin/sh

/etc/init.d/netdev_cfg.sh disable
/etc/init.d/mxl-network-init disable

exit 0

