# net_core/plugins/tr181-dhcpv4client/files/etc/uci-defaults/99-udhcpc-should-call-dhcpv4client-script
#!/bin/sh

sed -i "s/-s \/lib\/netifd\/dhcp.script/-s \/etc\/udhcpc.user/" /lib/netifd/proto/dhcp.sh

exit 0
