# net_services/plugins/deviceinfo-manager/files/etc/uci-defaults/99-prevent-uci-setting-hostname
#!/bin/sh

if [ -f "/usr/lib/amx/modules/mod-hostname-syscall.so" ]; then
    sed -i 's|echo "$hostname" > /proc/sys/kernel/hostname|#echo "$hostname" > /proc/sys/kernel/hostname|g' /etc/init.d/system
fi

exit 0
