# net_services/plugins/tr181-syslog/files/etc/uci-defaults/99-change-syslog-startup-order
#!/bin/sh

sed -i 's|START=.*|START=50|' /etc/init.d/syslog-ng

exit 0
