# net_core/plugins/dhcpv4-manager/files/etc/uci-defaults/99-remove-dnsmasq-hotplug-script
#!/bin/sh

if  [[ ! -f etc/hotplug.d/ntp/25-dnsmasqsec ]]; then
    rm -f /etc/hotplug.d/ntp/25-dnsmasqsec
fi

exit 0