# net_services/plugins/tr181-usermanagement/files/etc/uci-defaults/99-fix-default-users-passwd-field-in-passwd-and-shadow
#!/bin/sh

sed -i 's/^\([^:]*:\)[^:]*:/\1x:/' /etc/passwd
sed -i 's/^\([^:]*:\)[^:]*:/\1:/' /etc/shadow

exit 0
