# net_core/plugins/tr181-dns/files/etc/uci-defaults/99-dnsmasq-disable-dns
#!/bin/sh

sed -i "s/option port '53'/option port '0'/" /etc/config/dhcp

exit 0
