# prplos/network/plugins/tr181-syslog/files/etc/uci-defaults/99-change-syslog-startup-order
#!/bin/sh

sed -i 's|START=.*|START=11|' /etc/init.d/syslog-ng

exit 0
